# deploy/systemd/bbb-audit-publish.timer
[Unit]
Description=Publish audit data after the daily rollup

[Timer]
OnCalendar=*-*-* 03:10:00
Persistent=true
AccuracySec=1min
RandomizedDelaySec=30s
Unit=bbb-audit-publish.service

[Install]
WantedBy=timers.target
