# deploy/systemd/bbb-digest.service
[Unit]
Description=Bristol Bus Bot twice-daily estate digest

[Service]
Type=oneshot
User=@BBB_DEPLOY_USER@
Group=@BBB_DEPLOY_USER@
WorkingDirectory=@BBB_DEPLOY_BASE@/current/collector
EnvironmentFile=/etc/bristolbusbot/collector.env
Environment=BBB_BOT_DB=/var/lib/bristolbusbot/bot/app_data.db
ExecStart=/usr/local/libexec/bbb-run-recorded-job --name digest --skip-exit-code 75 -- /usr/bin/flock -n -E 75 /run/lock/bristolbusbot/digest.lock @BBB_DEPLOY_BASE@/current/collector/venv/bin/python @BBB_DEPLOY_BASE@/current/collector/status_digest.py
TimeoutStartSec=5min
UMask=0077
NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=/var/lib/bristolbusbot/monitoring /run/lock/bristolbusbot
CapabilityBoundingSet=
AmbientCapabilities=
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
