# deploy/systemd/bbb-staleness.timer
[Unit]
Description=Check collector freshness hourly

[Timer]
OnCalendar=*-*-* *:15:00
Persistent=true
AccuracySec=1min
RandomizedDelaySec=60s
Unit=bbb-staleness.service

[Install]
WantedBy=timers.target
