# deploy/systemd/bbb-backup.service
[Unit]
Description=Bristol Bus Bot encrypted local and R2 backup
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
Group=root
EnvironmentFile=/etc/bristolbusbot/backup.env
ExecStart=/usr/local/libexec/bbb-run-recorded-job --name backup -- /usr/local/libexec/bbb-backup-run
TimeoutStartSec=4h
UMask=0077
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
CPUWeight=20
NoNewPrivileges=yes
PrivateDevices=no
DevicePolicy=closed
PrivateTmp=yes
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=/mnt/bbb-backup /var/cache/bristolbusbot /var/lib/bristolbusbot/monitoring /run/lock
CapabilityBoundingSet=CAP_CHOWN CAP_DAC_OVERRIDE CAP_DAC_READ_SEARCH CAP_FOWNER
AmbientCapabilities=
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
